# files/etc/systemd/system/wg-quick@.service.d/reload-job.conf
[Service]
ExecReload=/bin/bash -c 'exec /usr/bin/wg syncconf %i <(exec /usr/bin/wg-quick strip %i)'
