# files/etc/systemd/system/wg-quick@wg0.service.d/pre-post-jobs.conf
[Service]
ExecStartPre=/opt/pivpn/network/networkCONF.sh wireguard apply
ExecStopPost=/opt/pivpn/network/networkCONF.sh wireguard reset
