# files/etc/systemd/system/openvpn-server@pivpn.service.d/pre-post-jobs.conf
[Service]
ExecStartPre=/opt/pivpn/network/networkCONF.sh openvpn apply
ExecStopPost=/opt/pivpn/network/networkCONF.sh openvpn reset
